# debian/conf/gunicorn_front.conf
[program:front]
user=encyc
directory=/usr/local/src/encyc-front/front
command=/usr/local/src/env/front/bin/gunicorn_django -w 2 -b 0.0.0.0:8080
autostart=true
autorestart=true
redirect_stderr=True
