# debian/etc/supervisor/conf.d/front.conf
[program:front]
command=/opt/front/frontenv/bin/gunicorn_django -w 2 -b 127.0.0.1:8002 /opt/front/settings.py
directory=/opt/front
user=www-data
autostart=true
autorestart=true
redirect_stderr=True
